# ids_project/ids-dashboard.service
[Unit]
Description=IDS Dashboard Server
After=network.target

[Service]
Type=simple
User=kali
Group=kali
WorkingDirectory=/home/kali/ids-project/ids_project
Environment="PATH=/home/kali/ids-project/ids_project/venv/bin"
ExecStart=/home/kali/ids-project/ids_project/venv/bin/python3 /home/kali/ids-project/ids_project/dashboard/start_server.py --server gunicorn --port 5000 --workers 4
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

